# packaging/systemd/gpt-model-proxy.service
[Unit]
Description=gpt-model-proxy
Documentation=https://github.com/BackFire/gpt-model-proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__GROUP__
ExecStart=__HOME__/.local/bin/gpt-model-proxy
Restart=on-failure
RestartSec=2s
Environment=HOME=__HOME__
Environment=PATH=__HOME__/.local/bin:/opt/gpt-model-proxy:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
